# com.nee2p.relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nee2p.relay</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/ttt/sites/2Pee/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ttt/sites/2Pee</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>9787</string>
        <key>HOST</key>
        <string>127.0.0.1</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ttt/Library/Logs/nee2p-relay.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ttt/Library/Logs/nee2p-relay.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
